# s6-rc.d/00-init-refresh/up
#!/bin/sh
exec /scripts/init-refresh.sh 2>&1