# etc/s6-overlay/s6-rc.d/git-daemon/run
#!/usr/bin/execlineb -P
fdmove -c 2 1
/etc/git-cache-server/scripts/exec.sh git daemon --base-path=/srv/git --verbose

