# s6-rc.d/spawn-fcgi/run
#!/usr/bin/execlineb -P
fdmove -c 2 1
spawn-fcgi -n -u www-data -s /var/run/fcgiwrap.socket /usr/sbin/fcgiwrap -- /usr/sbin/fcgiwrap -f
