# etc/s6-overlay/s6-rc.d/00-init-refresh/up
#!/bin/sh
exec /etc/git-cache-server/scripts/fetch.sh 2>&1