# etc/s6-overlay/s6-rc.d/supercronic/run
#!/usr/bin/execlineb -P
fdmove -c 2 1
/etc/git-cache-server/scripts/exec.sh /usr/local/bin/supercronic /etc/git-cache-server/crontab
